# armguard/deployment_A/methods/production/ci/systemd/armguard-staging.service
[Unit]
Description=ArmGuard staging (Gunicorn)
After=network.target postgresql.service redis.service
Requires=redis.service

[Service]
Type=simple
User=rds
Group=rds
WorkingDirectory=/opt/armguard
EnvironmentFile=/opt/armguard/.env.staging
Environment=DJANGO_SETTINGS_MODULE=core.settings_production
Environment=PATH=/opt/armguard/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/armguard/.venv/bin/gunicorn --workers 2 --bind 127.0.0.1:8010 core.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
